# deploy/systemd/cleanux-scheduler.timer.example
[Unit]
Description=Run CleanUx scheduler every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
Unit=cleanux-scheduler.service

[Install]
WantedBy=timers.target
