# deploy/supervisor/cleanux-reverb.conf
; CleanUx — Reverb WebSocket server
;
; À déployer dans /etc/supervisor/conf.d/cleanux-reverb.conf
; Puis :
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start cleanux-reverb:*
;
; Logs :
;   /var/log/supervisor/cleanux-reverb.log
;   /var/log/supervisor/cleanux-reverb-error.log

[program:cleanux-reverb]
process_name=%(program_name)s
command=php /var/www/cleanux/artisan reverb:start --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cleanux-reverb.log
stderr_logfile=/var/log/supervisor/cleanux-reverb-error.log
stopwaitsecs=10
; Augmenter la limite de fichiers ouverts (chaque connexion WebSocket = 1 fd).
; Cf. doc Reverb : https://laravel.com/docs/10.x/reverb#open-files
minfds=10000
; Si tu utilises ext-uv pour scaler au-delà de ~1000 connexions concurrentes,
; ajouter dans le command : --no-debug
