# deploy/supervisor/cleanux-worker.conf.example
[program:cleanux-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php8.5 /var/www/cleanux/artisan queue:work database --queue=default --sleep=3 --tries=3 --timeout=120 --max-time=3600
autostart=true
autorestart=true
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/cleanux/worker.log
stopwaitsecs=3600
